# djacha_supervisor.conf
[supervisord]
;nodaemon=true

[program:sshd]
directory=/home/ubuntu/djacha/
command=docker-compose run run_prod

